# installers/linux/voindns.service
[Unit]
Description=voindns local DoH proxy and DNS redirector
Documentation=https://voindns.net
After=network-online.target systemd-resolved.service
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/lib/voindns/voindns-service run
Restart=on-failure
RestartSec=3

# /run/voindns is created for the control socket.
RuntimeDirectory=voindns
RuntimeDirectoryMode=0755

# Bind :53 and edit system DNS. Runs as root for the MVP (systemd-resolved
# D-Bus + resolv.conf fallback need privilege). Tighten with a dedicated user +
# CAP_NET_ADMIN/CAP_NET_BIND_SERVICE in a follow-up.
Environment=VOINDNS_PORT=53
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
